# certbot-nginx/crontab
1 6,18 * * * sleep $(expr $RANDOM % 55)m; certbot renew --authenticator standalone --installer nginx --pre-hook "nginx -s stop"
